# moonclock-hardware.path
[Unit]
Description=Monitor database.json and restart moon-hardware

[Path]
PathModified=/usr/local/bin/moonclock/restart-hardware
Unit=moonclock-hardware.service

[Install]
WantedBy=default.target